# database/modules/elasticsearch/install.sh
#!/usr/bin/env bash

echo bla

wget https://download.elastic.co/elasticsearch/release/org/elasticsearch/distribution/deb/elasticsearch/2.3.1/elasticsearch-2.3.1.deb
sudo dpkg -i elasticsearch-2.3.1.deb
sudo systemctl enable elasticsearch.service
sudo rm elasticsearch-2.3.1.deb
